# templates/etcd.service.j2
[Unit]
Description=etcd key-value store
Documentation=https://github.com/coreos/etcd

[Service]
User=root
Group=root
Type=notify
ExecStart=/usr/bin/etcd{% if name %} --name {{ name }}{% endif %}{% if data_dir %} --data-dir {{ data_dir }}{% endif %}{% if wal_dir %} --wal-dir {{ wal_dir }}{% endif %}{% if snapshot_count %} --snapshot-count {{ snapshot_count }}{% endif %}{% if heartbeat_interval %} --heartbeat-interval {{ heartbeat_interval }}{% endif %}{% if election_timeout %} --election-timeout {{ election_timeout }}{% endif %}{% if initial_election_tick_advance %} --initial-election-tick-advance {{ initial_election_tick_advance }}{% endif %}{% if listen_peer_urls %} --listen-peer-urls {{ listen_peer_urls }}{% endif %}{% if listen_client_urls %} --listen-client-urls {{ listen_client_urls }}{% endif %}{% if max_snapshots %} --max-snapshots {{ max_snapshots }}{% endif %}{% if max_wals %} --max-wals {{ max_wals }}{% endif %}{% if cors %} --cors {{ cors }}{% endif %}{% if quota_backend_bytes %} --quota-backend-bytes {{ quota_backend_bytes }}{% endif %}{% if max_request_bytes %} --max-request-bytes {{ max_request_bytes }}{% endif %}{% if grpc_keepalive_min_time %} --grpc-keepalive-min-time {{ grpc_keepalive_min_time }}{% endif %}{% if grpc_keepalive_interval %} --grpc-keepalive-interval {{ grpc_keepalive_interval }}{% endif %}{% if grpc_keepalive_interval %} --grpc-keepalive-timeout {{ grpc_keepalive_interval }}{% endif %}{% if initial_advertise_peer_urls %} --initial-advertise-peer-urls {{ initial_advertise_peer_urls }}{% endif %}{% if initial_cluster %} --initial-cluster {{ initial_cluster }}{% endif %}{% if initial_cluster_state %} --initial-cluster-state {{ initial_cluster_state }}{% endif %}{% if initial_cluster_token %} --initial-cluster-token {{ initial_cluster_token }}{% endif %}{% if advertise_client_urls %} --advertise-client-urls {{ advertise_client_urls }}{% endif %}{% if discovery %} --discovery {{ discovery }}{% endif %}{% if discovery_fallback %} --discovery-fallback {{ discovery_fallback }}{% endif %}{% if discovery_proxy %} --discovery-proxy {{ discovery_proxy }}{% endif %}{% if discovery_srv %} --discovery-srv {{ discovery_srv }}{% endif %}{% if strict_reconfig_check %} --strict-reconfig-check {{ strict_reconfig_check }}{% endif %}{% if auto_compaction_retention %} --auto-compaction-retention {{ auto_compaction_retention }}{% endif %}{% if enable_v2 %} --enable-v2{% endif %}{% if proxy %} --proxy {{ proxy }}{% endif %}{% if proxy_failure_wait %} --proxy-failure-wait {{ proxy_failure_wait }}{% endif %}{% if proxy_failure_wait %} --proxy-refresh-interval {{ proxy_failure_wait }}{% endif %}{% if proxy_dial_timeout %} --proxy-dial-timeout {{ proxy_dial_timeout }}{% endif %}{% if proxy_write_timeout %} --proxy-write-timeout {{ proxy_write_timeout }}{% endif %}{% if proxy_read_timeout %} --proxy-read-timeout {{ proxy_read_timeout }}{% endif %}{% if cert_file %} --cert-file {{ cert_file }}{% endif %}{% if key_file %} --key-file {{ key_file }}{% endif %}{% if client_cert_auth %} --client-cert-auth{% elif not client_cert_auth %} --client-cert-auth 'false'{% endif %}{% if trusted_ca_file %} --trusted-ca-file {{ trusted_ca_file }}{% endif %}{% if auto_tls %} --auto-tls {{ auto_tls }}{% endif %}{% if peer_cert_file %} --peer-cert-file {{ peer_cert_file }}{% endif %}{% if peer_key_file %} --peer-key-file {{ peer_key_file }}{% endif %}{% if peer_client_cert_auth %} --peer-client-cert-auth{% elif not peer_client_cert_auth %} --peer-client-cert-auth 'false'{% endif %}{% if peer_trusted_ca_file %} --peer-trusted-ca-file {{ peer_trusted_ca_file }}{% endif %}{% if peer_auto_tls %} --peer-auto-tls {{ peer_auto_tls }}{% endif %}{% if cipher_suites %} --cipher-suites {{ cipher_suites }}{% endif %}{% if debug %} --debug {{ debug }}{% endif %}{% if proxy_read_timeout %} --log-package-levels {{ proxy_read_timeout }}{% endif %}{% if log_output %} --log-output {{ log_output }}{% endif %}{% if force_new_cluster %} --force-new-cluster {% endif %}{% if enable_pprof %} --enable-pprof {{ enable_pprof }}{% endif %}{% if metrics %} --metrics {{ metrics }}{% endif %}{% if auth_token %} --auth-token {{ auth_token }}{% endif %}

Restart=always
RestartSec=10s
LimitNOFILE=40000
TimeoutSec=30

[Install]
WantedBy=multi-user.target
